# systemd/pager-update.service
[Unit]
Description=Update Pager Software
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/pi/pagerpi/updates/update.sh

[Install]
WantedBy=pager.service